# configs/supervisord-spotweb.conf
[program:spotweb]
command=/start-spotweb-update.sh
startsecs=5
numprocs=1
autostart=true
autorestart=true